# robot_control/oled/oled-display.service
[Unit]
Description=Robot OLED WiFi Status Display
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory={{PROJECT_ROOT}}
ExecStart={{PROJECT_ROOT}}/.venv/bin/python {{PROJECT_ROOT}}/script/oled.py
Restart=always
RestartSec=5

# Ensure virtualenv and standard binaries are available.
Environment="PATH={{PROJECT_ROOT}}/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

[Install]
WantedBy=multi-user.target
